# overlay/etc/cron.d/root
# min   hour    day     month   weekday command
*       *       *       *       *       bf-cron -q 1min
*/15    *       *       *       *       bf-cron 15min
0       *       *       *       *       bf-cron hourly
0       2       *       *       *       bf-cron daily
0       3       *       *       6       bf-cron weekly
0       5       1       *       *       bf-cron monthly
