# overlay/etc/s6-overlay/s6-rc.d/init/up
foreground
{
    if { bf-echo "Initialising container." "init" }
    ifelse -n { with-contenv /usr/lib/bf/init }
    {
        if { bf-notok "Unable to initialise container." "init" }
        if { bf-cont-terminate }
    }
    if { bf-done "init" }
}
