# overlay/etc/s6-overlay/s6-rc.d/init/up
foreground
{
    ifelse -n { with-contenv bf-init }
    {
        if { bf-notok "Unable to initialise container." "init" }
        if { bf-cont-terminate }
    }
}
